# systemd/gardenctl.service.in
[Unit]
Description=Garden control MQTT client
Requires=network-online.target
After=sockets.target network-online.target

[Service]
Type=forking
PIDFile=/var/run/gardenctl.pid
ExecStart=/usr/local/bin/gardenctl -u "andreas" -s "drachen5Flug" --loglevel 2
ExecReload=/bin/kill -TERM $MAINPID
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
